# services/battery_monitor.service
[Unit]
Description=Battery Monitor Service
Documentation=Battery monitoring and shutdown service for UPS
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/sheepdog/Herd_OS
ExecStart=/usr/bin/python3 /home/sheepdog/Herd_OS/services/battery_monitor_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
